# configs/systemd/isp-rsc-backup.service.d/00-ispbilling-env.conf
[Service]
EnvironmentFile=-/etc/ispbilling.env
