# configs/systemd/isp-admin.service.d/scaling.conf
# __PHASE0_SCALING__ — added by phase01 script. Safe defaults to survive
# load spikes (auto-renew thundering herd, captive flood). Does NOT change
# the service behaviour day-to-day; only affects restart/kill thresholds.
[Service]
TimeoutStartSec=300
TimeoutStopSec=60
Restart=always
RestartSec=10
# Memory safety net — kills runaway workers BEFORE the OOM-killer takes
# the whole VPS down. 2G is generous for 4-worker uvicorn under normal load.
MemoryMax=2G
# Restart at most 5 times in 5 min — avoid death-spiral if config is broken.
StartLimitIntervalSec=300
StartLimitBurst=5
