# configs/systemd/vlan-radius-sync.service.d/00-venv.conf
[Service]
ExecStart=
ExecStart=/opt/ispbilling/venv/bin/python3 /opt/ispbilling/admin-portal/workers/vlan_radius_sync.py --interval=30
